# com.sshd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sshd.plist</string>
    <key>ProgramArguments</key>
    <array>
        <string>bash</string>
        <string>/Users/hannahcartier/.zshd/window_sync.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/hannahcartier/.zshd</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>out.log</string>
    <key>StandardErrorPath</key>
    <string>err.log</string>
</dict>
</plist>
